# yt-timestamps-watch.service
[Unit]
Description=Restart yt-timestamps bot when files change

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart yt-timestamps
